# server-root/etc/cron.d/sync
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

58 * * * * root sync; echo 3 > /proc/sys/vm/drop_caches
48 * * * * root fstrim /
